# bec-worker.service
[Unit]
Description=BEC Worker - %i
After=network-online.target
Wants=network-online.target

[Service]
Type=simple

User=ec2-user
Group=ec2-user

# Load environment variables (DB credentials, AWS config, etc.)
EnvironmentFile=/etc/bec/worker.env

# Working directory should be where the code is checked out
WorkingDirectory=/home/ec2-user/bec-orchestration

# Execute the worker using the installed CLI
# The %i specifier is the job name (e.g., ldv1) when using templated service
# For non-templated service, replace %i with your job name (e.g., ldv1)
ExecStart=/opt/pytorch/bin/bec worker --job-name %i

# Graceful shutdown handling for AWS ASG scale-down
# AWS sends SIGTERM and gives ~2 minutes (120s) before force termination
# Give worker 110s to complete current job gracefully
KillMode=mixed
KillSignal=SIGTERM
TimeoutStopSec=110

# Restart policy
Restart=always
RestartSec=5

# Send stdout/stderr to a log file (for CloudWatch Agent)
StandardOutput=append:/var/log/bec/worker-%i.log
StandardError=append:/var/log/bec/worker-%i.log

# Resource limits
LimitNOFILE=65536

# Security hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/var/log/bec /tmp

[Install]
WantedBy=multi-user.target
